# ansible/roles/containers/comfyui/templates/comfyui-idle-restart.service.j2
[Unit]
Description=Restart ComfyUI after idle generation
Wants={{ comfyui_service_name }}.service
After={{ comfyui_service_name }}.service docker.service

[Service]
Type=oneshot
ExecStart={{ comfyui_idle_restart_script }}
