# ansible/roles/headless_laptop/templates/asus-quiet-fan-policy.service.j2
[Unit]
Description=Apply ASUS quiet fan policy
ConditionPathExists=/sys/firmware/acpi/platform_profile

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/asus-quiet-fan-policy
